# systemd/overnight-coder.timer
[Unit]
Description=Schedule overnight coding sessions
Requires=overnight-coder.service

[Timer]
# Run at 11 PM every day
OnCalendar=*-*-* 23:00:00

# If we missed a run (system was off), run it when system starts
Persistent=true

# Add some randomness to avoid exact scheduling issues
RandomizedDelaySec=300

# Wake the system to run (if supported)
WakeSystem=false

[Install]
WantedBy=timers.target
